# states/neutron-network/restart_linuxbridge-agent.cron
#!/bin/bash

systemctl restart neutron-linuxbridge-agent

exit 0
